# watermark.service
[Unit]
Description=Watermark service
After=network.target

[Service]
Type=simple
User=ning
WorkingDirectory=/home/ning/src/watermark-demo
ExecStart=/bin/bash /home/ning/src/watermark-demo/run.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target